# 02_Setup a Load balancer/promisc-on-start.service
[Unit]
Description=Set promiscuous mode for wlan0 on startup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/ip link set wlan0 promisc on
ExecStop=/usr/bin/ip link set wlan0 promisc off
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target